# dist/argeo-node/rpm/usr/lib/systemd/system/argeo.service
[Unit]
Description=Argeo default node
After=network.target
Wants=postgresql.service

[Service]
Type=simple
ExecStart=/usr/sbin/argeoctl start
ExecReload=/usr/sbin/argeoctl reload
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
