# dist/argeo-node/rpm/usr/lib/systemd/system/argeo@.service
[Unit]
Description=Argeo Node %I
After=network.target
Wants=postgresql.service

[Service]
Type=simple
PIDFile=/var/run/argeo/%i.pid
ExecStart=/usr/sbin/argeoctl start %I
ExecReload=/usr/sbin/argeoctl reload %I
SuccessExitStatus=143
DynamicUser=true
StateDirectory=argeo.d/%I
LogsDirectory=argeo.d/%I
ConfigurationDirectory=argeo.d/%I

[Install]
WantedBy=multi-user.target
